# deploy/push-dispatch.service
# Oneshot service that fires the Web Push fan-out. Paired with
# push-dispatch.timer. Install alongside my-koto.service.
[Unit]
Description=my-koto Web Push dispatch
After=my-koto.service
Wants=my-koto.service

[Service]
Type=oneshot
# Run dir holds the runtime env; the dispatch script lives in the git clone.
WorkingDirectory=/srv/my-koto
EnvironmentFile=/srv/my-koto/.env.production.local
ExecStart=/opt/koto-city/deploy/push-dispatch.sh
